# amdfand.service
[Unit]
Description=amdfan controller
[Service]
ExecStart=/usr/bin/amdfand service
Restart=always
[Install]
WantedBy=multi-user.target
