# services/amdvold
#!/sbin/openrc-run

description="amdvol controller"
pidfile="/run/${SVCNAME}.pid"
command="/usr/bin/amdvold"
command_args="service"
command_user="root"

depend() {
    need udev
}
